# etc/rc.local
#!/bin/bash

# issue with 5850u 
# disables cpu threads
echo 'off' | sudo tee /sys/devices/system/cpu/smt/control